# realm/initramfs/etc/inittab
# Setup filesystem
::sysinit:/bin/busybox --install
::sysinit:/bin/mkdir -p /proc
::sysinit:/bin/mkdir -p /var/run
::sysinit:/bin/mkdir -p /etc/network/if-pre-up.d
::sysinit:/bin/mkdir -p /etc/network/if-down.d
::sysinit:/bin/mkdir -p /etc/network/if-up.d
::sysinit:/bin/mkdir -p /etc/network/if-post-down.d
::sysinit:/bin/mkdir -p /dev
::sysinit:/bin/mkdir -p /sys
::sysinit:/bin/mount -t devtmpfs devtmpfs /dev
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys

# Run init scripts
/dev/console::sysinit:/etc/init.d/rcS

# Spawn user console
/dev/console::askfirst:-/bin/sh

# On restart
::restart:/init

# Shutdown on Ctrl+Alt+Del
::ctrlaltdel:/bin/poweroff -f


